# supervisor.conf
[supervisord]
nodaemon=true
user=root


[program:grpc]
command=bash start.sh
user=root
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:http]
command=bash web.sh
user=root
